# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/books15/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/books15/crontab.log 2>&1
